# fabtools/tests/cron.d/broken_cron

# START test (generated by fabtools)
0 */2 * * * username /home/username/test.pl
# END test (generated by fabtools)


# START test (generated by fabtools)
0 */2 * * * username /home/username/test.pl
# END test (generated by fabtools)


# START test2 (generated by fabtools)
0 */2 * * * username /home/username/test.pl
# END test2 (generated by fabtools)


# START test3 (generated by fabtools)
0 */2 * * * username /home/username/test.pl

# END test2 (generated by fabtools)
